# checkoff.timer
[Unit]
Description=Poll GW2 API for account updates every 15 mins

[Timer]
OnCalendar=*:5/15

[Install]
WantedBy=default.target
